# debate-bot/debate-bot.service
[Unit]
Description=Telegram Debate Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=sreekar
WorkingDirectory=/var/www/html/debate-bot
EnvironmentFile=/var/www/html/debate-bot/.env
ExecStart=/var/www/html/debate-bot/venv/bin/python bot.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
