# debate-bot/debate-bot.service
[Unit]
Description=Telegram Debate Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/debate-bot
EnvironmentFile=/home/pi/debate-bot/.env
ExecStart=/home/pi/debate-bot/venv/bin/python bot.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
